# pin-hwmon-metrics.service
[Unit]
Description=pin-hwmon metrics exporter

[Service]
Type=oneshot
ExecStart=/usr/bin/pin-hwmon-metrics